# deploy/evoiceclaw-desktop-v3.service
# ──────────────────────────────────────────────────────────────
# eVoiceClaw Desktop v3 — systemd 服务单元文件
#
# 安装步骤：
#   1. 将此文件复制到远程服务器：
#      scp deploy/evoiceclaw-desktop-v3.service <user>@<host>:/etc/systemd/system/
#
#   2. 重新加载 systemd 配置：
#      systemctl daemon-reload
#
#   3. 启用并启动服务：
#      systemctl enable evoiceclaw-desktop-v3
#      systemctl start evoiceclaw-desktop-v3
#
#   4. 查看状态和日志：
#      systemctl status evoiceclaw-desktop-v3
#      journalctl -u evoiceclaw-desktop-v3 -f
# ──────────────────────────────────────────────────────────────

[Unit]
# 服务描述
Description=eVoiceClaw Desktop v3 — AI OS 后端服务
# 网络就绪后再启动
After=network.target
# 如果使用了外部数据库，可取消下面的注释
# After=network.target postgresql.service

[Service]
# 服务类型：简单前台进程
Type=simple

# 运行用户和组（根据实际情况修改）
User=root
Group=root

# 工作目录
WorkingDirectory=/opt/evoiceclaw-desktop-v3

# 启动命令
ExecStart=/usr/bin/env uvicorn app.main:app --host 0.0.0.0 --port 28771 --workers 1

# 优雅停止：先发 SIGTERM，等待 30 秒后发 SIGKILL
ExecStop=/bin/kill -SIGTERM $MAINPID
TimeoutStopSec=30

# 重启策略：异常退出时自动重启
Restart=on-failure
# 重启间隔（秒）
RestartSec=5
# 最大连续重启次数（60 秒内最多重启 3 次）
StartLimitBurst=3
StartLimitIntervalSec=60

# 环境变量
Environment=PYTHONUNBUFFERED=1
Environment=LOG_LEVEL=INFO

# 如果使用虚拟环境，取消注释并修改路径
# Environment=PATH=/opt/evoiceclaw-desktop-v3/.venv/bin:/usr/local/bin:/usr/bin:/bin
# ExecStart=/opt/evoiceclaw-desktop-v3/.venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 28771 --workers 1

# 标准输出和错误输出写入 journal
StandardOutput=journal
StandardError=journal
# journal 日志标识符
SyslogIdentifier=evoiceclaw-desktop-v3

# 安全加固（可选，根据需求启用）
# 禁止提升权限
NoNewPrivileges=true
# 保护系统目录（只读）
ProtectSystem=strict
# 保护 /home 目录
ProtectHome=true
# 允许写入的目录（工作目录下的 data 和 logs）
ReadWritePaths=/opt/evoiceclaw-desktop-v3/data /opt/evoiceclaw-desktop-v3/logs /opt/evoiceclaw-desktop-v3
# 私有临时目录
PrivateTmp=true

[Install]
# 多用户模式下自动启动
WantedBy=multi-user.target
